# etc/update-motd.d/10-weather
#!/bin/sh
echo
echo "CURRENT WEATHER AT LINZ"
ansiweather -l linz -u metric